# contrib/etc/init/happo-agent-autoscaling-leave.conf
description "happo-agent leave at node shutdown"
author  "Takatada Yoshima <yoshima@heartbeats.jp>"

start on runlevel [016]

task

exec /usr/local/bin/happo-agent leave 2>&1 | logger -t happo-agent
